# ai-security-proxy.service
[Unit]
Description=AI Security Proxy Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=www-data
WorkingDirectory=/opt/ai-security-proxy
ExecStart=/opt/ai-security-proxy/run_proxy_prod.sh --bind 0.0.0.0 --rate-limit 60
ExecStop=/bin/kill -s TERM $MAINPID

# Environment variables
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=0
Environment=PROXY_PORT=8080
Environment=HEALTH_PORT=8081
Environment=ENABLE_HEALTH_CHECK=true
Environment=ENABLE_RATE_LIMITING=true
Environment=BASIC_AUTH_ENABLED=true
Environment=LOG_LEVEL=info

# Security settings
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
MemoryDenyWriteExecute=true

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096
CPUQuota=80%
MemoryLimit=2G

[Install]
WantedBy=multi-user.target